# blinkie/crontab
* * * * * root node /blinkie/resetRateLimits.js
